# mnt/mtd/config/crontabs
1 3 * * * /mnt/mtd/netvue/firmware/bin/upgrade_check.sh
